# app-admin/iptsd/files/iptsd.initd
#!/sbin/openrc-run
# Copyright 1999-2020 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

description="Intel Precise Touch & Stylus daemon"
command=/usr/bin/iptsd
command_args=""
pidfile="/run/${RC_SVCNAME}.pid"

# Let start-stop-daemon handle backgrounding and the pidfile
command_background="true"
